# etc/cron.d/fastpinger
# Run fastpinger every minute
* * * * * source ~/etc/fastchecker.conf; $BASE_DIR/fastpinger;
